# rf-monitor/rf-ha-status.service
[Unit]
Description=Publish RF Monitor status to Home Assistant MQTT
After=network-online.target rf-anomaly-detector.service
Wants=network-online.target

[Service]
Type=oneshot
User=cmilkosk
ExecStart=/usr/bin/python3 /home/cmilkosk/rf-monitor/publish_rf_ha_status.py
