# hackrf-monitor/hackrf-influx.service
[Unit]
Description=HackRF sweep to InfluxDB collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cmilkosk
EnvironmentFile=/home/cmilkosk/.config/hackrf-influx.env
ExecStart=/usr/bin/python3 /home/cmilkosk/bin/hackrf_influx_collector.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
